# systemd/karton-system@.service.d/10-setup-bucket.conf
[Service]
Environment=EXTRA_ARGS=--setup-bucket
