# systemd/karton-logger@.service.d/10-command.conf
[Service]
ExecStart=
ExecStart=/opt/extractor/bin/karton --config-file /opt/extractor/etc/extractor.ini logs
