# systemd/karton-logger@.service.d/10-command.conf
[Service]
ExecStart=
ExecStart=/opt/expander/bin/karton --config-file /opt/expander/etc/expander.ini logs
